# systemd/user/backup.service
[Unit]
Description=full system backup

[Service]
ExecStart=%h/bin/backup.sh
