# systemd/user/bupBackup.service
[Unit]
Description=full system backup with bup

[Service]
ExecStart=${HOME}/bin/backup.sh
